# payphone-app/wyoming-whisper.service
[Unit]
Description=Wyoming Hailo Whisper STT Server
Documentation=https://github.com/CarbonRaven/payphone-ai
After=network.target
Before=payphone.service

[Service]
User=tom
WorkingDirectory=/home/tom/tele-ai/payphone-app
ExecStart=/home/tom/tele-ai/payphone-app/.venv/bin/python \
    services/wyoming_whisper_server.py \
    --port 10300 \
    --hef /usr/local/hailo/resources/models/hailo10h/Whisper-Base.hef \
    --model-dir /home/tom/tele-ai/payphone-app/models \
    --variant base \
    --log-level INFO
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
